# examples/python-runtime-sandbox/supervisord.conf
# Supervisord config: 与 openclaw 官方 Dockerfile 一致（见 repo Dockerfile 末行 CMD），再叠加 python-sandbox
# 官方 CMD 为 node openclaw.mjs gateway --allow-unconfigured；ps 里 openclaw / openclaw-gateway 为同一进程不同阶段的 process.title
# Run with: supervisord -n -c ./supervisord.conf (from container WORKDIR)
# Socket 放在 /tmp，便于以非 root (runAsUser: 1000) 运行时 supervisorctl 能连上

[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# 与官方 Dockerfile CMD 一致：node openclaw.mjs gateway；--bind lan --port 18789 便于从宿主机/Ingress 访问
[program:openclaw-gateway]
command=node openclaw.mjs gateway --allow-unconfigured --bind lan --port 18789
directory=/app
environment=HOME="/home/node"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:python-sandbox]
command=python3 -m uvicorn main:app --host 0.0.0.0 --port 8888 --log-level trace
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
